# docker-cm-data-collector/rootfs/etc/supervisor.d/data-collector.ini

[program:data-collector]
command        = /usr/local/bin/data-collector.rb
stdout_logfile = /tmp/%(program_name)s.log
stderr_logfile = /tmp/%(program_name)s.log
autorestart    = true
