# docker-cm-smashing/rootfs/etc/supervisor.d/dashing.ini

[program:dashing]
directory      = /opt/%DASHBOARD%
command        = /usr/bin/thin --pid /tmp/thin.pid --log /tmp/thin.log --port 3030 --chdir /opt/%DASHBOARD% --environment development --rackup /opt/%DASHBOARD%/config.ru start
stdout_logfile = /tmp/%(program_name)s.log
stderr_logfile = /tmp/%(program_name)s.log
autorestart    = true
