# monitoring/docker-graphite/rootfs/etc/supervisor.d/memcached.ini

[program:memcached]
command        = /usr/bin/memcached -l 127.0.0.1 -m 4 -u memcached
stdout_logfile = /tmp/%(program_name)s.log
stderr_logfile = /tmp/%(program_name)s.log
autorestart    = true
