# environments/aio/share/carbon-relay/etc/supervisor.d/carbon-relay.ini
[program:carbon-relay]
command        = /usr/bin/carbon-c-relay -f /etc/carbon-c-relay.conf -w 4 -q 2056 -m -b 2056 -l /var/log/carbon-relay.log
stdout_logfile = /tmp/%(program_name)s.log
stderr_logfile = /tmp/%(program_name)s.log
autorestart    = true
